# roles/borg/templates/backup.timer.j2
[Unit]
Description=Executes BORG Backup {{ borg_name }}

[Timer]
OnCalendar={{ borg_schedule }}

[Install]
WantedBy=timers.target
