# roles/borg/templates/backup.service.j2
[Unit]
Description=BORG Backup {{ borg_name }}

[Service]
Type=oneshot
ExecStart=/etc/borg/{{ borg_name }}_run.sh
